# crontab
@reboot screen -S anki-autoupdate "/Users/$USER/Desktop/Git/personal-mnemonic-medium/update_mnemonic_medium.sh"
